# deploy/github-deploy-webhook.service
# Canonical production owner for the authenticated GitHub deploy webhook.
# The repair workflow delegates all service installation to this unit's installer.
[Unit]
Description=Metrotherapy GitHub deploy webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/metrotherapy
EnvironmentFile=/etc/metrotherapy/github-deploy-webhook.env
ExecStart=/root/metrotherapy/.venv/bin/python /root/deploy_webhook.py
Restart=always
RestartSec=3

# The webhook only authenticates one fixed endpoint and asks systemd to start a
# fixed deploy worker. Keep the root process heavily sandboxed even though the
# transient deploy unit itself still needs full deployment privileges.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077

[Install]
WantedBy=multi-user.target
