# deploy/github-deploy-webhook.service
# Canonical production owner for the authenticated GitHub deploy webhook.
# The repair workflow delegates all service installation to this unit's installer.
[Unit]
Description=Metrotherapy GitHub deploy webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/metrotherapy
EnvironmentFile=/etc/metrotherapy/github-deploy-webhook.env
ExecStart=/root/metrotherapy/.venv/bin/python /root/deploy_webhook.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
