# deploy/metrotherapy.service
[Unit]
Description=Metrotherapy Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=metrotherapy
Group=metrotherapy
WorkingDirectory=/opt/metrotherapy
EnvironmentFile=-/etc/metrotherapy/metrotherapy.env
EnvironmentFile=-/etc/default/metrotherapy
Environment=APP_ENV=prod
Environment=PYTHONUNBUFFERED=1
Environment=HEALTHCHECK_ENABLED=1
Environment=HEALTHCHECK_HOST=127.0.0.1
Environment=HEALTHCHECK_PORT=8082
Environment=MESSENGER_WEBHOOK_HOST=127.0.0.1
Environment=MESSENGER_WEBHOOK_PORT=8081
Environment=METRO_DB_ENGINE=sqlite
Environment=DATABASE_URL=
Environment=METRO_DB_PATH=/var/lib/metrotherapy/data.db
Environment=LOG_PATH=/var/log/metrotherapy/app.log
ExecStart=/opt/metrotherapy/.venv/bin/python main.py
Restart=always
RestartSec=3
TimeoutStartSec=45
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictNamespaces=true
UMask=0077
StateDirectory=metrotherapy
LogsDirectory=metrotherapy
ReadWritePaths=/var/lib/metrotherapy /var/log/metrotherapy
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
